# cron/crontab
5 0 * * * curl -X POST "http://meal-voucher:8000/meal-plan-recurring/apply?api_key=${CRON_API_KEY}"
